# aprox-clean.service.in
[Unit]
Description=APT Proxy Clean-up

[Service]
CacheDirectory=apt/proxy
DynamicUser=yes
ExecStart=@BINDIR@/aprox-clean
