# aprox.service.in
[Unit]
Description=APT Proxy
Requires=aprox.socket

[Service]
CacheDirectory=apt/proxy
DynamicUser=yes
#Environment=CACHE_DIR=/var/cache/apt/proxy
ExecStart=@BINDIR@/aprox
