# aprox.service.in
[Unit]
Description=APT Proxy
Requires=aprox.socket

[Service]
CacheDirectory=apt/proxy
DynamicUser=yes
ExecStart=@BINDIR@/aprox
